# scripts/.scripts/shell/system/udev/rules.d/95-monitor.rules
KERNEL=="card0", SUBSYSTEM=="drm", ENV{DISPLAY}=":0", ENV{XAUTHORITY}="/home/lady/.Xauthority", RUN+="/usr/bin/systemctl start monitor.service"
